# deploy/oracle/systemd/t-employ-health.service
[Unit]
Description=TBELL Employ — ops health check
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/t-employ
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/t-employ
Nice=5
TimeoutStartSec=120
ExecStart=/usr/local/bin/t-employ-run ops:health
