# deploy/oracle/systemd/t-employ-session.timer
[Unit]
Description=TBELL Employ — daily session refresh 06:30 KST
Requires=t-employ-session.service

[Timer]
OnCalendar=*-*-* 06:30:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-session.service

[Install]
WantedBy=timers.target
