# deploy/oracle/systemd/t-employ-session.timer
[Unit]
Description=TBELL Employ — session refresh 06:47 & 13:17 KST (weekdays)
Requires=t-employ-session.service

[Timer]
OnCalendar=Mon..Fri *-*-* 06:47:00
OnCalendar=Mon..Fri *-*-* 13:17:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-session.service

[Install]
WantedBy=timers.target
