# deploy/oracle/systemd/t-employ-digest.timer
[Unit]
Description=TBELL Employ — morning digest Mon-Fri 07:30 KST
Requires=t-employ-digest.service

[Timer]
OnCalendar=Mon..Fri *-*-* 07:30:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-digest.service

[Install]
WantedBy=timers.target
