# deploy/oracle/systemd/t-employ-digest.service
[Unit]
Description=TBELL Employ — morning digest mail
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/default/t-employ
User=ubuntu
Group=ubuntu
WorkingDirectory=/opt/t-employ
Nice=5
TimeoutStartSec=1800
ExecStart=/usr/local/bin/t-employ-run mail:morning-digest
