# deploy/oracle/systemd/t-employ-talent.timer
[Unit]
Description=TBELL Employ — daily talent crawl 07:00 KST
Requires=t-employ-talent.service

[Timer]
OnCalendar=*-*-* 07:00:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-talent.service

[Install]
WantedBy=timers.target
