# deploy/oracle/systemd/t-employ-pdf.timer
[Unit]
Description=TBELL Employ — applicant PDF batch 07:20 KST
Requires=t-employ-pdf.service

[Timer]
OnCalendar=*-*-* 07:20:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-pdf.service

[Install]
WantedBy=timers.target
