# deploy/oracle/systemd/t-employ-health.timer
[Unit]
Description=TBELL Employ — health check every 6h
Requires=t-employ-health.service

[Timer]
OnCalendar=*-*-* 00/6:10:00
Timezone=Asia/Seoul
Persistent=true
Unit=t-employ-health.service

[Install]
WantedBy=timers.target
